# deployment/aws/systemd.service
# Example systemd unit when running the app without Docker (venv + uvicorn).
# Copy to /etc/systemd/system/ai-log-analyzer.service and adjust paths/user.

[Unit]
Description=AI Log Analyzer
After=network.target

[Service]
Type=simple
User=appuser
WorkingDirectory=/opt/ai-log-analyzer
EnvironmentFile=/opt/ai-log-analyzer/.env
ExecStart=/opt/ai-log-analyzer/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
